# storageagent/storageagent.service
[Unit]
Description=KVM CSI Driver StorageAgent

[Service]
ExecStart=/usr/local/bin/storageagent

[Install]
WantedBy=multi-user.target
